# snapraid-sync.service
[Unit]
Description=Run SnapRAID sync command
 
[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/bin/snapraid sync
